# autostart/autokey.desktop
[Desktop Entry]
Type=Application
Name=Autokey
Exec=/usr/bin/autokey
Terminal=false
