# autostart/copyq.desktop
[Desktop Entry]
Type=Application
Name=CopyQ
Exec=/usr/bin/copyq
Terminal=false
