# autostart/vscode.desktop
[Desktop Entry]
Type=Application
Name=VSCode
Exec=/usr/bin/code
Terminal=false
