# autostart/vscode-snap.desktop
[Desktop Entry]
Type=Application
Name=VSCode
Exec=/snap/bin/code
Terminal=false
